# .config/systemd/user/aro-update.service
[Unit]
Description=Update aro

[Service]
WorkingDirectory=%h/aro
ExecStart=%h/aro/update.sh
